# gearmanworker.ini
[program:gearmanworker]
environment=GEARMAN_USER="root",PATH="/usr/bin":%(ENV_PATH)s
command=/bin/bash -c "/bin/sleep 2 && php /root/t3census-worker-detector/worker/Typo3HostDetectorWorker.php"
directory=/root/t3census-worker-detector
exitcodes=0
stdout_logfile=/gearman/worker_stdout.log
stdout_logfile_maxbytes=1MB
sterr_logfile=/gearman/worker_stderr.log
stderr_logfile_maxbytes=1MB
redirect_stderr=true
process_name=%(program_name)s_%(process_num)02d
numprocs_start=0
numprocs=50
startsecs=0
autostart=true
autorestart=true
startretries=5
user=root
stopsignal=KILL
